# deploy/mod-retention/qiki-mod-retention.service
[Unit]
Description=Qikiworld moderation evidence cleanup
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/qikiworld/server
ExecStart=/usr/bin/node dist/src/jobs/modRetention.js
