# deploy/mod-retention/qiki-mod-retention.timer
[Unit]
Description=Run moderation evidence cleanup daily

[Timer]
OnCalendar=*-*-* 03:10:00
Persistent=true
RandomizedDelaySec=60

[Install]
WantedBy=timers.target
